# etc/service/php-fpm/run
#!/bin/bash -e

# Is there any other script to run here?
[ -f /etc/service/php-fpm/run.initialization ] && source /etc/service/php-fpm/run.initialization

exec /etc/init.d/php-fpm start 2>&1
